# launchd/com.vibesquad.chrome.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Install template: replace __VAULT_ROOT__ (and __HOME__ if present) with absolute paths before launchctl load. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibesquad.chrome</string>
    <key>ProgramArguments</key>
    <array>
        <string>__VAULT_ROOT__/bin/chrome-bootstrap.sh</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/vibesquad-chrome-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/vibesquad-chrome-stderr.log</string>
</dict>
</plist>
